# systemd/uvicorn.service
[Unit]
Description=Uvicorn ASGI server for Aleqsys.com
After=network.target

[Service]
Type=notify
User=rovak
Group=webgroup
RuntimeDirectory=gunicorn
WorkingDirectory=/var/www/aleqsys.com/aleqsys_app
Environment="PATH=/var/www/aleqsys.com/.venv/bin"
Environment="DJANGO_SETTINGS_MODULE=aleqsys_app.settings.production"
EnvironmentFile=/home/rovak/aleqsys_app/.env

# Uvicorn with Gunicorn workers for production
ExecStart=/var/www/aleqsys.com/.venv/bin/gunicorn aleqsys_app.asgi:application \
    -k uvicorn.workers.UvicornWorker \
    --workers 3 \
    --bind unix:/run/gunicorn/gunicorn.sock \
    --access-logfile /var/log/aleqsys/uvicorn-access.log \
    --error-logfile /var/log/aleqsys/uvicorn-error.log \
    --capture-output \
    --enable-stdio-inheritance

ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
